# pve-ramfs-manager.service
[Unit]
Description=Proxmox VE RAMFS Manager - Reduces disk I/O by storing PVE data in RAM
Documentation=https://github.com/AndreiMaelstrom/pve-ramfs-manager
Before=pve-cluster.service rrdcached.service
After=local-fs.target
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=notify
ExecStart=/usr/local/bin/pve-ramfs-manager.sh start
ExecStop=/usr/local/bin/pve-ramfs-manager.sh stop
NotifyAccess=all
RemainAfterExit=yes
TimeoutStartSec=60
TimeoutStopSec=120

# Ensure data is persisted before shutdown
KillMode=process
KillSignal=SIGTERM

# Resource limits
MemoryMax=1G
TasksMax=10

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/pve-cluster /var/lib/pve-cluster-persistent /var/lib/rrdcached/db /var/lib/rrdcached-persistent /dev/shm /var/log
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
RequiredBy=pve-cluster.service rrdcached.service
