# pve-ramfs-manager.service
[Unit]
Description=Proxmox VE RAMFS Manager - Reduces disk I/O by storing PVE data in RAM
Documentation=https://github.com/AndreiMaelstrom/pve-ramfs-manager
Before=pve-cluster.service rrdcached.service
After=local-fs.target
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=notify
ExecStart=/usr/local/bin/pve-ramfs-manager.sh start
NotifyAccess=all
TimeoutStartSec=60
TimeoutStopSec=120

# Ensure data is persisted before shutdown
KillMode=process
KillSignal=SIGTERM
Restart=on-failure
RestartSec=5s

# Resource limits
MemoryMax=1G
TasksMax=10

# Runtime directory for lock file
RuntimeDirectory=pve-ramfs-manager

NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
# Note: Before= in [Unit] section ensures proper ordering with pve-cluster and rrdcached
# We don't use RequiredBy= to avoid blocking those services if this fails
